# recipe-book.service
[Unit]
Description=Recipe Book Web Application
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/recipe-book
Environment=FLASK_ENV=production
ExecStart=/usr/bin/python3 /home/pi/recipe-book/app.py
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/recipe-book/app.log
StandardError=append:/var/log/recipe-book/error.log

[Install]
WantedBy=multi-user.target
